# installer/linux/zarafa-dagent.init.alt
#!/bin/bash
#
# zarafa-dagent Zarafa Collaboration Platform's Delivery Agent
#
# chkconfig: - 86 24
# description: The Zarafa Delivery Agent in LMTP mode can be used to \
#              run the zarafa-dagent as a daemon. The Zarafa Delivery \
#              Agent can also be used as a standalone program.
# processname: /usr/bin/zarafa-dagent
# config: /etc/zarafa/dagent.cfg
# pidfile: /var/run/zarafa-dagent.pid

### BEGIN INIT INFO
# Provides: zarafa-dagent
# Required-Start: $local_fs $network $remote_fs $syslog
# Required-Stop: $local_fs $network $remote_fs $syslog
# Should-Start: zarafa-server
# Should-Stop: zarafa-server
# Short-Description: Zarafa Collaboration Platform's Delivery Agent
# Description:       The Zarafa Delivery Agent in LMTP mode can be used to
#                    run the zarafa-dagent as a daemon. The Zarafa Delivery
#                    Agent can also be used as a standalone program.
### END INIT INFO

# Do not load RH compatibility interface
WITHOUT_RC_COMPAT=1

# Source function library
. /etc/init.d/functions

DAGENTCONFIG=/etc/zarafa/dagent.cfg
DAGENTPROGRAM=/usr/bin/zarafa-dagent
dagent=`basename $DAGENTPROGRAM`
PIDFILE=/var/run/$dagent.pid
LOCKFILE=/var/lock/subsys/$dagent
RETVAL=0

# Sanity checks.
[ -x $DAGENTPROGRAM ] || exit 0

# the -d option is required to run LMTP mode deamonized
DAGENTCONFIG_OPT="-d"
[ ! -z $DAGENTCONFIG -a -f $DAGENTCONFIG ] && DAGENTCONFIG_OPT="$DAGENTCONFIG_OPT -c $DAGENTCONFIG"

[ -f /etc/sysconfig/zarafa ] && . /etc/sysconfig/zarafa
if [ -z "$ZARAFA_LOCALE" ]; then
	ZARAFA_LOCALE="C"
fi

start() {
	start_daemon --pidfile "$PIDFILE" --lockfile "$LOCKFILE" $dagent $DAGENTCONFIG_OPT
	RETVAL=$?
	
	return $RETVAL
}

stop() {
	stop_daemon --pidfile "$PIDFILE" --lockfile "$LOCKFILE" $dagent
	RETVAL=$?
	
	return $RETVAL
}

restart() {
	stop
	start
}

reload() {
	msg_reloading $dagent
	stop_daemon --pidfile "$PIDFILE" -HUP $dagent
	RETVAL=$?

	return $RETVAL
}

# See how we were called.
case "$1" in
    start)
		start
		;;
    stop)
		stop
		;;
    reload)
		reload
		;;
    restart|force-reload)
		restart
		;;
    condstop)
		if [ -e "$LOCKFILE" ]; then
		stop
		fi
		;;
    condrestart)
		if [ -e "$LOCKFILE" ]; then
		restart
		fi
		;;
    condreload)
		if [ -e "$LOCKFILE" ]; then
		reload
		fi
		;;
    status)
		status --pidfile "$PIDFILE" $dagent
		RETVAL=$?
		;;
    *)
		msg_usage "${0##*/} {start|stop|reload|restart|condstop|condrestart|condreload|status}"
		RETVAL=1
esac

exit $RETVAL
