# ECtools/zarafa-stats/zarafa-ssm.cron
# every sunday, 00:05
5 0 * * 0 root [ -x /usr/bin/zarafa-ssm ] && /usr/bin/zarafa-ssm
